# plugins/fancontrol/scripts/install
#!/bin/bash
# Aquacomputer Octo Fan Control Setup Script
# Configures fancontrol for CPU temperature-based fan curves
# Compatible with fresh OS installs - auto-detects hwmon devices

set -e

echo "=== Aquacomputer Octo Fan Control Setup ==="
echo

# Check if running as root
if [ "$EUID" -ne 0 ]; then
    echo "Please run as root (use sudo)"
    exit 1
fi

# Install fancontrol if not present
if ! command -v fancontrol &> /dev/null; then
    echo "Installing lm-sensors (fancontrol)..."
    if command -v dnf &> /dev/null; then
        dnf install -y lm_sensors
    elif command -v apt &> /dev/null; then
        apt install -y lm-sensors
    elif command -v pacman &> /dev/null; then
        pacman -S --noconfirm lm_sensors
    else
        echo "ERROR: Unknown package manager. Please install lm-sensors manually."
        exit 1
    fi
fi

echo "✓ fancontrol is installed"
echo

# Auto-detect hwmon devices
echo "Detecting hardware monitoring devices..."

# Find k10temp (CPU temperature sensor)
K10TEMP_HWMON=""
for hwmon in /sys/class/hwmon/hwmon*; do
    if [ -f "$hwmon/name" ] && [ "$(cat $hwmon/name)" = "k10temp" ]; then
        K10TEMP_HWMON=$(basename $hwmon)
        break
    fi
done

# Find Octo (fan controller)
OCTO_HWMON=""
for hwmon in /sys/class/hwmon/hwmon*; do
    if [ -f "$hwmon/name" ] && [ "$(cat $hwmon/name)" = "octo" ]; then
        OCTO_HWMON=$(basename $hwmon)
        break
    fi
done

# Validate detection
if [ -z "$K10TEMP_HWMON" ]; then
    echo "ERROR: Could not find k10temp sensor"
    echo "Make sure your CPU temperature sensor is available"
    exit 1
fi

if [ -z "$OCTO_HWMON" ]; then
    echo "ERROR: Could not find Aquacomputer Octo"
    echo "Make sure the Octo is connected and the aquacomputer_d5next driver is loaded"
    exit 1
fi

echo "✓ Found k10temp at $K10TEMP_HWMON"
echo "✓ Found Octo at $OCTO_HWMON"
echo

# Get device paths for config
K10TEMP_DEVPATH=$(readlink -f /sys/class/hwmon/$K10TEMP_HWMON | sed 's|/sys/||' | sed "s|/hwmon/$K10TEMP_HWMON||")
OCTO_DEVPATH=$(readlink -f /sys/class/hwmon/$OCTO_HWMON | sed 's|/sys/||' | sed "s|/hwmon/$OCTO_HWMON||")

echo "Generating fancontrol configuration..."

# Create fancontrol config
cat > /etc/fancontrol << EOF
# Fancontrol configuration for Aquacomputer Octo
# Auto-generated by setup-fancontrol.sh
# CPU temp source: k10temp ($K10TEMP_HWMON)
# Fan controller: Octo ($OCTO_HWMON)

INTERVAL=2

# Device paths (auto-detected)
DEVPATH=$K10TEMP_HWMON=$K10TEMP_DEVPATH $OCTO_HWMON=$OCTO_DEVPATH
DEVNAME=$K10TEMP_HWMON=k10temp $OCTO_HWMON=octo

# Temperature source for each PWM (all use CPU temp)
FCTEMPS=$OCTO_HWMON/pwm1=$K10TEMP_HWMON/temp1_input $OCTO_HWMON/pwm2=$K10TEMP_HWMON/temp1_input $OCTO_HWMON/pwm3=$K10TEMP_HWMON/temp1_input $OCTO_HWMON/pwm4=$K10TEMP_HWMON/temp1_input $OCTO_HWMON/pwm5=$K10TEMP_HWMON/temp1_input $OCTO_HWMON/pwm6=$K10TEMP_HWMON/temp1_input $OCTO_HWMON/pwm7=$K10TEMP_HWMON/temp1_input $OCTO_HWMON/pwm8=$K10TEMP_HWMON/temp1_input

# Fan monitoring (each PWM monitors its corresponding fan)
FCFANS=$OCTO_HWMON/pwm1=$OCTO_HWMON/fan1_input $OCTO_HWMON/pwm2=$OCTO_HWMON/fan2_input $OCTO_HWMON/pwm3=$OCTO_HWMON/fan3_input $OCTO_HWMON/pwm4=$OCTO_HWMON/fan4_input $OCTO_HWMON/pwm5=$OCTO_HWMON/fan5_input $OCTO_HWMON/pwm6=$OCTO_HWMON/fan6_input $OCTO_HWMON/pwm7=$OCTO_HWMON/fan7_input $OCTO_HWMON/pwm8=$OCTO_HWMON/fan8_input

# Case fans (1-6): QUIET curve 40°C-80°C (prioritizes silence)
# At 40°C: 15% speed (very quiet, barely audible)
# At 60°C: 50% speed (still quiet)
# At 80°C: 100% speed (max cooling when needed)
MINTEMP=$OCTO_HWMON/pwm1=40 $OCTO_HWMON/pwm2=40 $OCTO_HWMON/pwm3=40 $OCTO_HWMON/pwm4=40 $OCTO_HWMON/pwm5=40 $OCTO_HWMON/pwm6=40
MAXTEMP=$OCTO_HWMON/pwm1=80 $OCTO_HWMON/pwm2=80 $OCTO_HWMON/pwm3=80 $OCTO_HWMON/pwm4=80 $OCTO_HWMON/pwm5=80 $OCTO_HWMON/pwm6=80

# CPU fans (7-8): QUIET but responsive curve 45°C-78°C
# At 45°C: 15% speed (whisper quiet at idle)
# At 61.5°C: 54% speed (moderate)
# At 78°C: 100% speed (full cooling under sustained load)
MINTEMP=$OCTO_HWMON/pwm7=45 $OCTO_HWMON/pwm8=45
MAXTEMP=$OCTO_HWMON/pwm7=78 $OCTO_HWMON/pwm8=78

# Minimum PWM values (15% = 38/255, quietest safe speed)
MINSTART=$OCTO_HWMON/pwm1=38 $OCTO_HWMON/pwm2=38 $OCTO_HWMON/pwm3=38 $OCTO_HWMON/pwm4=38 $OCTO_HWMON/pwm5=38 $OCTO_HWMON/pwm6=38 $OCTO_HWMON/pwm7=38 $OCTO_HWMON/pwm8=38
MINSTOP=$OCTO_HWMON/pwm1=38 $OCTO_HWMON/pwm2=38 $OCTO_HWMON/pwm3=38 $OCTO_HWMON/pwm4=38 $OCTO_HWMON/pwm5=38 $OCTO_HWMON/pwm6=38 $OCTO_HWMON/pwm7=38 $OCTO_HWMON/pwm8=38
MINPWM=$OCTO_HWMON/pwm1=38 $OCTO_HWMON/pwm2=38 $OCTO_HWMON/pwm3=38 $OCTO_HWMON/pwm4=38 $OCTO_HWMON/pwm5=38 $OCTO_HWMON/pwm6=38 $OCTO_HWMON/pwm7=38 $OCTO_HWMON/pwm8=38

# Maximum PWM (100% = 255)
MAXPWM=$OCTO_HWMON/pwm1=255 $OCTO_HWMON/pwm2=255 $OCTO_HWMON/pwm3=255 $OCTO_HWMON/pwm4=255 $OCTO_HWMON/pwm5=255 $OCTO_HWMON/pwm6=255 $OCTO_HWMON/pwm7=255 $OCTO_HWMON/pwm8=255
EOF

echo "✓ Configuration written to /etc/fancontrol"
echo

# Enable and start service
echo "Enabling fancontrol service..."
systemctl enable fancontrol
systemctl restart fancontrol

echo
echo "=== Setup Complete ==="
echo
echo "Fancontrol is now running and will start automatically at boot."
echo
echo "Current status:"
systemctl status fancontrol --no-pager | head -10
echo
echo "To customize fan curves, edit /etc/fancontrol and restart the service:"
echo "  sudo systemctl restart fancontrol"
echo
echo "Fan curve summary (QUIET profile):"
echo "  Case fans (1-6): 40°C → 15% | 60°C → 50% | 80°C → 100%"
echo "  CPU fans (7-8):  45°C → 15% | 61.5°C → 54% | 78°C → 100%"
echo
echo "Safety note: Modern CPUs are safe up to 95°C. These curves prioritize"
echo "silence while maintaining excellent thermal performance."
